# attic/boxes/blog/files/etc/systemd/system/set-hostname.service
[Unit]
Description=Set the hostname on boot

[Service]
ExecStart=/bin/bash -c 'hostname blog'

[Install]
WantedBy=multi-user.target
