# boxes/attic/git/files/etc/systemd/system/set-hostname.service
[Unit]
Description=Set the hostname on boot

[Service]
ExecStart=/bin/bash -c 'hostname git'

[Install]
WantedBy=multi-user.target
